# recipes-tinker/heatmaster/heatmaster/heatmaster.service
[Unit]
Description=HeatMaster
After=multi-user.target

[Service]
ExecStart=/usr/bin/heatmaster
Restart=always

[Install]
WantedBy=multi-user.target
